# systemd/rbackupd@.timer
[Unit]
Description="Periodic restic backups"

[Timer]
OnBootSec=2m
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
