# systemd/rbackupd@.timer
[Unit]
Description="Periodic restic backups"

[Timer]
OnBootSec=2m
OnUnitInactiveSec=1d
Persistent=true

[Install]
WantedBy=default.target
